# pinnacle-security.service
[Unit]
Description=Pinnacle Copilot WindSurf Cascade Sync Service
After=network.target

[Service]
User=root
WorkingDirectory=/Users/samoey/Desktop/PinnacleCopilot
ExecStart=/usr/bin/python3 /Users/samoey/Desktop/PinnacleCopilot/sync_manager.py
Restart=always
RestartSec=30
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=pinnacle_sync

[Install]
WantedBy=multi-user.target
